# upstart.conf
#!upstart
# Fichier de configuration upstart
description "mes-infos-nutritionnelles"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec "/root/projects/mes-infos-nutritionnelles/run.sh" >> /var/log/mes-infos-nutritionnelles.log 2>&1

